# files/etc/uci-defaults/98-tinyfilemanager
#!/bin/sh
unzip -o /root/tinyfilemanager.zip -d /www
/etc/init.d/tinyfilemanager reload
rm -rf /root/tinyfilemanager.zip