# files/etc/uci-defaults/cloud_flare
#!/bin/sh
cp -r /root/cloudflared /usr/bin/
chmod +x /usr/bin/cloudflared
/etc/init.d/cloudflared reload
rm -rf /root/cloudflared