# files/etc/uci-defaults/tinyfilemanager
#!/bin/sh
tar -xf /root/tinyfilemanager.zip -C /www/
/etc/init.d/tinyfilemanager reload
rm -rf /root/tinyfilemanager.zip