# files/etc/uci-defaults/tinyfilemanager
#!/bin/sh
cp -r /root/tinyfilemanager/ /www/
/etc/init.d/tinyfilemanager reload
rm -rf /root/tinyfilemanager